# debian/etc/systemd/system/audio-cd-rip.service
[Unit]
Description=Headless CD Ripper

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/audio-cd-rip.sh
User=root
