# debian/etc/systemd/system/movie-dvd-rip.service
[Unit]
Description=Headless DVD Ripper

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/movie-dvd-rip.sh
User=root
